# magic_dingus_box_cpp/systemd/magic-dingus-box-cpp.service
[Unit]
Description=Magic Dingus Box C++ Kiosk Engine
After=network-online.target systemd-user-sessions.service sound.target
Wants=network-online.target sound.target
# Start early, before X11 (we use DRM/KMS directly)
Before=graphical.target lightdm.service
# Don't require X11 - we bypass it entirely
Conflicts=lightdm.service

[Service]
Type=simple
# Run as magic user (member of video group)
User=magic
Group=magic
# Working directory where executable and assets are located
# Adjust this path to match your installation location
WorkingDirectory=/opt/magic_dingus_box/magic_dingus_box_cpp/build
# Wait for HDMI audio and start PulseAudio
ExecStartPre=/bin/bash /opt/magic_dingus_box/magic_dingus_box_cpp/scripts/init_audio.sh
ExecStart=/opt/magic_dingus_box/magic_dingus_box_cpp/build/magic_dingus_box_cpp
ExecStopPost=/usr/bin/pulseaudio --kill
# Don't restart automatically - we handle restarts manually for RetroArch integration
# Restart=on-failure
# RestartSec=5
# Give service time to start properly on boot
StartLimitIntervalSec=300
StartLimitBurst=5
# Standard output/error logging
StandardOutput=journal
StandardError=journal
# Environment (if needed)
Environment=HOME=/home/magic
Environment=DISPLAY=
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Kill all processes in the service's cgroup on stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
# Start at multi-user target (before X11)
WantedBy=multi-user.target
