# systemd/magic-mpv.service
[Unit]
Description=Magic Dingus mpv
After=sound.target graphical.target lightdm.service
Wants=lightdm.service
# Wait for X11 socket and user session to be ready
After=graphical-session.target

[Service]
Type=simple
User=alexanderchaney
Group=alexanderchaney
RuntimeDirectory=magic
Environment=AUDIO_DEVICE=alsa/plughw:CARD=vc4hdmi0,DEV=0
EnvironmentFile=-/etc/magic_dingus_box.env
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/alexanderchaney/.Xauthority
# Wait for X11 socket and XAUTHORITY to be ready (max 30 seconds)
ExecStartPre=/bin/bash -c 'for i in {1..30}; do [ -S /tmp/.X11-unix/X0 ] && [ -f "$HOME/.Xauthority" ] && break || sleep 1; done; sleep 2'
ExecStart=/usr/bin/mpv --idle=yes --no-osc --no-osd-bar --keep-open=yes \
  --no-config \
  --input-ipc-server=/run/magic/mpv.sock \
  --input-vo-keyboard=no --input-default-bindings=no \
  --vo=x11 --hwdec=v4l2m2m \
  --video-sync=desync --video-latency-hacks=no \
  --vd-lavc-threads=4 --vd-lavc-fast=yes \
  --vd-lavc-skiploopfilter=all --vd-lavc-skipframe=nonref \
  --vd-lavc-dr=yes --sws-fast=yes \
  --cache=yes --cache-secs=30 --demuxer-max-bytes=500M --demuxer-max-back-bytes=500M \
  --force-window=yes --no-border \
  --loop-playlist=no --loop-file=no --loop=no
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target

