# systemd/magic-mpv-x11.service
[Unit]
Description=Magic Dingus mpv (X11)
After=sound.target graphical.target lightdm.service
Wants=lightdm.service
After=graphical-session.target

[Service]
Type=simple
User=magic
Group=magic
RuntimeDirectory=magic
Environment=AUDIO_DEVICE=alsa/plughw:CARD=vc4hdmi0,DEV=0
EnvironmentFile=-/etc/magic_dingus_box.env
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/magic/.Xauthority
# Set CPU governor to performance mode for better video playback
# Temporarily disabled - may cause UI rendering issues
# ExecStartPre=/opt/magic_dingus_box/scripts/set_cpu_performance.sh
# Wait for X11 socket and XAUTHORITY to be ready (max 30 seconds)
ExecStartPre=/bin/bash -c 'for i in {1..30}; do [ -S /tmp/.X11-unix/X0 ] && [ -f "$HOME/.Xauthority" ] && break || sleep 1; done; sleep 2'
ExecStart=/opt/magic_dingus_box/scripts/mpv_with_bezel.sh
Restart=on-failure

[Install]
WantedBy=graphical.target

