# systemd/magic-ui-home.service
[Unit]
Description=Magic Dingus UI (Home Directory)
After=graphical.target lightdm.service magic-mpv.service
Wants=magic-mpv.service
Requires=magic-mpv.service

[Service]
Type=simple
User=alexanderchaney
Group=alexanderchaney
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/var/run/lightdm/alexanderchaney/:0
WorkingDirectory=/home/alexanderchaney/magic_dingus_box
# Wait for X server to be ready (simple delay)
ExecStartPre=/bin/sleep 10
ExecStart=/home/alexanderchaney/magic_dingus_box/venv/bin/python -m magic_dingus_box.main
Restart=always
RestartSec=5
# Don't restart too quickly if it fails
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=graphical.target

