# systemd/magic-ui-x11.service
[Unit]
Description=Magic Dingus UI (X11)
After=network-online.target graphical.target lightdm.service
Wants=network-online.target lightdm.service
After=graphical-session.target
# Ensure system is fully ready before starting
After=systemd-user-sessions.service
# Note: magic-mpv.service is a system service, so we can't require it directly
# The ExecStartPre script will wait for the mpv socket to be ready

[Service]
Type=simple
# User and Group are not needed for user services - they run as the logged-in user
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/alexanderchaney/.Xauthority
Environment=SDL_AUDIODRIVER=alsa
Environment=MPV_SOCKET=/run/magic/mpv.sock
WorkingDirectory=/opt/magic_dingus_box
# Clean up stale RetroArch lock file on boot (if process is not running)
ExecStartPre=/bin/bash -c 'LOCK_FILE="/tmp/magic_retroarch_active.lock"; if [ -f "$LOCK_FILE" ]; then LOCK_PID=$(cat "$LOCK_FILE" 2>/dev/null | tr -d "[:space:]"); if [ -n "$LOCK_PID" ] && [ "$LOCK_PID" != "starting" ]; then if ! ps -p "$LOCK_PID" >/dev/null 2>&1 || ! ps -p "$LOCK_PID" -o comm= 2>/dev/null | grep -qi retroarch; then rm -f "$LOCK_FILE" 2>/dev/null; fi; else rm -f "$LOCK_FILE" 2>/dev/null; fi; fi'
# Wait for X11 socket and XAUTHORITY to be ready (max 60 seconds for boot)
# On boot, things take longer to initialize, so we wait up to 60 seconds
# Note: mpv socket will be created by mpv service, but we don't wait for it here
# The app will handle mpv connection retries internally
ExecStartPre=/bin/bash -c 'for i in {1..60}; do [ -S /tmp/.X11-unix/X0 ] && [ -f "$HOME/.Xauthority" ] && break || sleep 1; done; sleep 3'
ExecStart=/opt/magic_dingus_box/venv/bin/python -m magic_dingus_box.main
Restart=on-failure
RestartSec=5
# Give service time to start properly on boot
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target

