# systemd/magic-ui.service
[Unit]
Description=Magic Dingus UI
After=network-online.target graphical.target lightdm.service
Wants=network-online.target lightdm.service
# Note: mpv is started internally by the UI application
# Wait for X11 socket and user session to be ready
After=graphical-session.target

[Service]
Type=simple
User=alexanderchaney
Group=alexanderchaney
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XAUTHORITY=/home/alexanderchaney/.Xauthority
WorkingDirectory=/opt/magic_dingus_box
# Clean up stale RetroArch lock file on boot (if process is not running)
ExecStartPre=/bin/bash -c 'LOCK_FILE="/tmp/magic_retroarch_active.lock"; if [ -f "$LOCK_FILE" ]; then LOCK_PID=$(cat "$LOCK_FILE" 2>/dev/null | tr -d "[:space:]"); if [ -n "$LOCK_PID" ] && [ "$LOCK_PID" != "starting" ]; then if ! ps -p "$LOCK_PID" >/dev/null 2>&1 || ! ps -p "$LOCK_PID" -o comm= 2>/dev/null | grep -qi retroarch; then rm -f "$LOCK_FILE" 2>/dev/null; fi; else rm -f "$LOCK_FILE" 2>/dev/null; fi; fi'
# Wait for X11 socket and XAUTHORITY to be ready (max 30 seconds)
ExecStartPre=/bin/bash -c 'for i in {1..30}; do [ -S /tmp/.X11-unix/X0 ] && [ -f "$HOME/.Xauthority" ] && break || sleep 1; done; sleep 2'
ExecStart=/usr/bin/python3 -m magic_dingus_box.main
Restart=always
RestartSec=5

[Install]
# Requires X11 for --wid embedding
WantedBy=graphical.target

